# docker/cfg/supervisor/conf.d/app.conf
[program:backend]
directory=/opt/waterfall/backend
command=uwsgi /etc/uwsgi/backend.ini -L
redirect_stderr=true
